# nuc-assistant.service
[Unit]
Description=NUC Offline Voice Assistant
After=network.target ollama.service sound.target
Wants=ollama.service
Documentation=https://github.com/<your-username>/<your-repo>

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=__REPO_PATH__
ExecStartPre=/bin/sleep 5
ExecStart=__PYTHON_BIN__ __REPO_PATH__/assistant.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nuc-assistant

# Audio access
SupplementaryGroups=audio

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=PULSE_RUNTIME_PATH=__PULSE_RUNTIME_PATH__

[Install]
WantedBy=multi-user.target
